# examples/crontab
* * * * * root php /var/www/piplin/artisan schedule:run >> /dev/null 2>&1
